# root/cyclictest/overlay/etc/inittab

# linux-conf-perf init script
console::sysinit:/usr/bin/linux-conf-perf -L console 115200 vt100
